# analysis/device_binaries_full/usr/lib/lib/ipkg/info/atfwd-daemon.list
/etc/init.d/start_atfwd_daemon
/usr/bin/atfwd_daemon
