# analysis/device_binaries_full/usr/lib/lib/ipkg/info/avahi.list
/usr/lib/libavahi-common.so
/usr/sbin/avahi-daemon
/usr/lib/libavahi-core.so.7
/etc/rc.d/S35avahi-daemon
/etc/init.d/avahi-daemon
/usr/lib/libavahi-core.so
/usr/share/avahi/service-types
/usr/share/avahi/avahi-service.dtd
/etc/avahi/hosts
/usr/lib/libavahi-core.so.7.0.2
/etc/avahi/avahi-daemon.conf
/usr/lib/libavahi-common.so.3.5.3
/usr/lib/libavahi-common.so.3
