# analysis/device_binaries_full/usr/lib/lib/ipkg/info/data-ipa-cfg-mgr.list
/usr/lib/libipanat.so.1
/data/misc/ipa/IPACM_cfg.xml
/usr/bin/ipacm
/usr/lib/libipanat.so.1.0.0
/etc/init.d/start_ipacm_le
