# analysis/device_binaries_full/usr/lib/lib/ipkg/info/data.list
/WEBSERVER/www/QCMAP_WWAN.html
/usr/bin/QCMAP_ConnectionManager
/usr/bin/qti_socksv5
/usr/lib/libqcmap_cm.so.1.0.0
/usr/lib/libqcmaputils.so.1
/usr/lib/libdsi_netctrl.so.0
/usr/bin/ipacm_perf
/etc/udhcpc.d/udhcpc.script
/WEBSERVER/www/QCMAP.html
/WEBSERVER/www/QCMAP_Firewall_Help.html
/usr/lib/libqdi.so.0
/etc/init.d/port_bridge
/usr/lib/libdsi_netctrl.so.0.0.0
/etc/data/netmgr_config.xml
/WEBSERVER/www/QCMAP_Account_Help.html
/usr/lib/libqmi_ip.so.1
/etc/init.d/start_ipacm_perf_le
/WEBSERVER/www/cgi-bin/qcmap_auth
/etc/data/dsi_config.xml
/data/qti_socksv5_auth.xml
/etc/init.d/start_QCMAP_ConnectionManager_le
/etc/init.d/start_qti_le
/WEBSERVER/www/QCMAP_IPV6_Help.html
/usr/bin/QCMAP_StaInterface
/WEBSERVER/www/QCMAP_WLAN.html
/etc/init.d/start_eMBMs_TunnelingModule_le
/usr/lib/libqcmap_cm.so.1
/usr/bin/qmi_ip_multiclient
/usr/lib/libqdp.so.0
/usr/bin/port_bridge
/usr/bin/QCMAP_Web_CLIENT
/usr/lib/libqdi.so.0.0.0
/usr/lib/libqcmap_client.so.1.0.0
/data/l2tp_cfg.xml
/WEBSERVER/www/QCMAP_WLAN_Help.html
/WEBSERVER/www/QCMAP_Firewall.html
/etc/init.d/start_ipacmdiag_le
/WEBSERVER/www/QCMAP_Account.html
/usr/lib/libnetmgr.so.0
/usr/lib/libqdp.so.0.0.0
/usr/lib/libqmi_ip.so.1.0.0
/WEBSERVER/www/QCMAP_UPNP_DLNA_MDNS_Help.html
/usr/lib/libqcmap_client.so.1
/etc/init.d/netmgrd
/usr/bin/qti_ppp
/data/qti_socksv5_conf.xml
/WEBSERVER/www/QCMAP_NAT_Help.html
/usr/bin/eMBMs_TunnelingModule
/usr/lib/libqcmapipc.so.1.0.0
/etc/init.d/data-init
/etc/qmi_ip_cfg.xml
/usr/lib/libqcmapipc.so.1
/WEBSERVER/www/QCMAP_IPV4_Help.html
/WEBSERVER/www/QCMAP_WWAN_Help.html
/WEBSERVER/www/QCMAP_IPV4.html
/usr/tests/dsi_netctrl_test
/usr/lib/libqcmaputils.so.1.0.0
/etc/init.d/start_stop_qmi_ip_multiclient
/WEBSERVER/www/cgi-bin/qcmap_web_cgi
/WEBSERVER/www/QCMAP_NAT.html
/usr/bin/ipacmdiag
/data/qmi_ip_cfg.xml
/WEBSERVER/www/QCMAP_UPNP_DLNA_MDNS.html
/usr/bin/QCMAP_CLI
/usr/bin/netmgrd
/WEBSERVER/www/images/delete.png
/usr/lib/libnetmgr.so.0.0.0
/usr/bin/qti
/WEBSERVER/www/QCMAP_Help.html
/WEBSERVER/www/images/header.jpg
/WEBSERVER/www/images/update.jpg
/data/mobileap_cfg.xml
/WEBSERVER/www/QCMAP_login.html
/etc/init.d/start_stop_qti_ppp_le
/WEBSERVER/www/QCMAP_IPV6.html
/usr/bin/radish
