# analysis/device_binaries_full/usr/lib/lib/ipkg/info/diag.list
/usr/bin/diag_klog
/usr/bin/test_diag
/usr/bin/diag_dci_sample
/etc/init.d/chgrp-diag
/usr/lib/libdiag.so.1.0.0
/usr/bin/diag_socket_log
/usr/lib/libdiag.so.1
/usr/bin/PktRspTest
/usr/bin/diag_uart_log
/usr/bin/diag_mdlog
/usr/bin/diag_callback_sample
