# analysis/device_binaries_full/usr/lib/lib/ipkg/info/inetutils.list
/usr/bin/rcp
/etc/init.d/syslogd
/usr/bin/rsh
/etc/rcS.d/S20syslogd
/etc/syslog.conf
/sbin/syslogd
