# analysis/device_binaries_full/usr/lib/lib/ipkg/info/init-mss.list
/etc/init.d/init_sys_mss
/sbin/init_mss
