# analysis/device_binaries_full/usr/lib/lib/ipkg/info/powerapp.list
/data/reboot-cookie
/sbin/powerapp
/etc/init.d/reset_reboot_cookie
/sbin/reboot-bootloader
/sbin/sys_reboot
/sbin/reboot-recovery
/sbin/sys_shutdown
