# analysis/device_binaries_full/usr/lib/lib/ipkg/info/samba.list
/etc/rc.d/S95samba
/usr/private/secrets.tdb
/usr/lib/libtalloc.so.2
/usr/lib/libtevent.so
/usr/sbin/smbd
/etc/samba/smb.conf
/usr/lib/libwbclient.so.0
/etc/init.d/samba
/usr/bin/smbstatus
/usr/bin/smbcontrol
/usr/bin/nmblookup
/usr/bin/pdbedit
/usr/bin/tdbbackup
/usr/lib/libtdb.so.1
/usr/lib/libtalloc.so.2.0.5
/usr/lib/libtdb.so.1.2.9
/usr/bin/testparm
/usr/lib/samba/lowcase.dat
/usr/lib/samba/valid.dat
/usr/lib/samba/upcase.dat
/usr/bin/net
/usr/lib/libtevent.so.0
/usr/lib/libtdb.so
/usr/lib/libwbclient.so
/usr/lib/samba/vfs/audit.so
/usr/sbin/nmbd
/usr/lib/samba/vfs/full_audit.so
/usr/lib/libtalloc.so
/usr/lib/libtevent.so.0.9.11
