# analysis/device_binaries_full/usr/lib/lib/ipkg/info/start-scripts-firmware-links.list
/lib/firmware/image
/etc/init.d/firmware-links.sh
