# analysis/device_binaries_full/usr/lib/lib/ipkg/info/start-scripts-ipa-fws.preinst
#!/bin/sh
if [ -z "$D" -a -f "/etc/init.d/ipa_fws.sh" ]; then
	/etc/init.d/ipa_fws.sh stop
fi
if type update-rc.d >/dev/null 2>/dev/null; then
	if [ -n "$D" ]; then
		OPT="-f -r $D"
	else
		OPT="-f"
	fi
	update-rc.d $OPT ipa_fws.sh remove
fi
