# analysis/device_binaries_full/usr/lib/lib/ipkg/info/system-core-adbd.list
/etc/launch_adbd
/usr/lib/libadbd.so.0.0.0
/usr/lib/libadbd.so.0
/sbin/adbd
/etc/init.d/adbd
