# analysis/device_binaries_full/usr/lib/lib/ipkg/info/tftp-server.list
/system/rfs/msm/mpss/hlos
/system/rfs/mdm/mpss/readonly/firmware
/system/rfs/msm/mpss/ramdumps
/system/rfs/msm/mpss/readonly/firmware
/system/rfs/msm/adsp/readwrite
/system/rfs/mdm/adsp/readwrite
/system/rfs/mdm/mpss/hlos
/system/rfs/mdm/adsp/shared
/system/rfs/mdm/mpss/shared
/system/rfs/mdm/adsp/ramdumps
/system/rfs/msm/adsp/shared
/system/rfs/msm/mpss/shared
/etc/init.d/tftp_server.sh
/system/rfs/msm/mpss/readwrite
/system/rfs/mdm/mpss/readwrite
/system/rfs/mdm/mpss/ramdumps
/system/rfs/msm/adsp/hlos
/system/rfs/msm/adsp/ramdumps
/system/rfs/mdm/adsp/hlos
/system/rfs/mdm/adsp/readonly/firmware
/system/rfs/msm/adsp/readonly/firmware
/sbin/tftp_server
