# analysis/device_binaries_full/usr/lib/lib/ipkg/info/udev.list
/lib/udev/rules.d/95-udev-late.rules
/lib/udev/rules.d/60-persistent-serial.rules
/lib/udev/udevd
/bin/udevadm
/etc/init.d/udev
/etc/rc.d/S00udev
/lib/udev/rules.d/75-tty-description.rules
/lib/libudev.so.0.13.1
/etc/udev/udev.conf
/lib/libudev.so
/lib/libudev.so.0
