# analysis/device_binaries_full/usr/lib/lib/ipkg/info/vsftpd.list
/etc/vsftpd.chroot_list
/etc/rc3.d/K20vsftpd3
/etc/vsftpd.conf
/etc/rc2.d/K20vsftpd2
/etc/rc4.d/S85vsftpd4
/etc/init.d/vsftpd
/etc/rc1.d/K20vsftpd1
/etc/vsftpd.user_list
/usr/bin/vsftpd
